# templates/app.service.j2
[Unit]
Description=Lab 3 ansible
Requires=network.target
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=java -jar {{ app_dir }}/{{ app_file }}
Environment=PORT={{ app_port }}